# build/provisioning/roles/backend/files/etc/systemd/system/skedbackend.service
[Unit]
Description=SkedBacked service
After=network.target

[Service]
Type=simple
Restart=on-failure
User=skedbackend
ExecStart=/opt/skedbackend/skedbackend
StandardOutput=append:/var/log/skedbackend/app.log
StandardError=append:/var/log/skedbackend/error.log
WorkingDirectory=/opt/skedbackend

[Install]
WantedBy=multi-user.target
