# build/provisioning/roles/backend/templates/etc/systemd/system/skedbackend.service.d/override.conf.j2
[Service]
Environment="PG_STR=postgres://sked:{{ pg_passwd }}@localhost:54320/sked?sslmode=disable"
Environment="SENDGRID_API_KEY={{ sg_api_key }}"
Environment="FROM_EMAIL={{ from_email }}"
