# build/provisioning/roles/backend/files/templates/etc/systemd/system/skedbackend.service.d/override.conf.j2
[Service]
Environment="PG_STR="
Environment="SENDGRID_API_KEY="
Environment="FROM_EMAIL="
